# systemd/.config/systemd/user/bluetuih.service
[Unit]
Description=Starts bluetuith in the backgound

[Service]
ExecStart=/sbin/urxvt -e bluetuith
Restart=always

[Install]
WantedBy=default.target
