# systemd/.config/systemd/user/chatgpt_chat.service
[Unit]
Description=Starts chatgpt-desktop in the backgound

[Service]
ExecStart=/sbin/chat-gpt
Restart=always

[Install]
WantedBy=default.target
