# systemd/.config/systemd/user/ncpamixer.service
[Unit]
Description=Starts ncpamixer in the backgound

[Service]
ExecStart=/sbin/kitty -e ncpamixer
Restart=always

[Install]
WantedBy=default.target
