# systemd/.config/systemd/user/snixembed.service
[Unit]
Description=Starts sixembed used to detect icons(needed for talon icon)

[Service]
ExecStart=/sbin/snixembed
Restart=always

[Install]
WantedBy=default.target
